# deploy/habitat-api.service
[Unit]
Description=Habitat API and dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=habitat
Group=habitat
WorkingDirectory=/opt/habitat-cli
EnvironmentFile=/etc/habitat/habitat-api.env
ExecStart=/home/habitat/.bun/bin/bun /opt/habitat-cli/src/server.ts
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
